# backup/vcs-noshell.in.sv.gnu.org/etc/cron.d/truncate-gitcvs-db-log
MAILTO=savannah-reports-private@gnu.org
0 3 9 * * nobody : > /srv/git/gnulib.git/gitcvs-db/log
