# backup/sftp.in.sv.gnu.org/etc/cron.d/download-timestamp
*/15 * * * *	root	(LANG=C; date +\%s; date -u) > /srv/download/00_TIME.txt
