# backup/colonialone.fsf.org/etc/cron.daily/backup-bind
#!/bin/bash
rsync -qvHS --numeric-ids --delete root@savannah-backup.gnu.org:/etc/bind/ /root/20_backup-etc-bind/
